# linux/config/systemd/user/mbsync.timer
[Unit]
Description=Mysync synchronization timer

[Timer]
OnBootSec=1m
OnUnitActiveSec=20m
Unit=mbsync.service

[Install]
WantedBy=timers.target
