# linux/config/systemd/user/wallpaper.timer
[Unit]
Description=Runs bg.sh every day

[Timer]
OnBootSec=0min
OnUnitActiveSec=691min

[Install]
WantedBy=timers.target
