# linux/config/systemd/user/xremap.service
# systemctl --user add-wants niri.service xremap.service
[Unit]
Description=Run xremap 
# BindsTo=sway-session.target
# # After=network.target
# After=udevmon.service
# DefaultDependencies=no
# PartOf=graphical-session.target
# After=graphical-session.target
# Requisite=graphical-session.target
After=default.target

[Service]
# Environment="LANG=zh_CN.UTF-8"
# Environment="XDG_SESSION_DESKTOP=niri"
Type=simple
# User=jixiuf
# Group=users
Restart=always
# KillMode=process prevents systemd from killing child processes (e.g., emacsclient)
# This avoids "No child processes (ECHILD)" errors when xremap spawns subprocesses
KillMode=process
# Separate dbus environment setup to reduce process nesting and improve subprocess handling
ExecStartPre=/usr/bin/dbus-update-activation-environment --systemd --all
# sh -c is needed to expand ${XDG_SESSION_DESKTOP};
# exec replaces shell with xremap so systemd tracks the correct PID
ExecStart=/bin/sh -c 'exec /usr/local/bin/xremap-${XDG_SESSION_DESKTOP} %h/.config/xremap/xremap-${XDG_SESSION_DESKTOP}.yaml %h/.config/xremap/xremap.yaml  --watch=device --ignore="dotool keyboard" --ignore="Yubico YubiKey OTP+CCID"'

[Install]
# WantedBy=sway-session.target
WantedBy=default.target

# Local Variables:
# eval: (add-hook (make-local-variable 'after-save-hook) #'(lambda()(shell-command "systemctl --user daemon-reload;systemctl --user restart xremap")) t)
# mode: conf
# End:
