# linux/etc/systemd:system:emerge-sync.service
[Unit]
Description=Sync Portage Tree

[Service]
Type=oneshot
ExecStart=/usr/bin/emerge --sync
[Install]
WantedBy=multi-user.target
