# linux/etc/systemd:system:emerge-sync.timer
[Unit]
Description=Daily Portage Tree Sync

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
