# linux/etc/systemd:system:rsync.service
[Unit]
Description=Backs up files

[Service]
ExecStart=/usr/bin/rsync -avz -CERrltm --delete /var/cache/binpkgs root@10.2.0.171:/data/binpkgs 
